# remind/.config/systemd/user/hb.service
[Unit]
Description=Birthday Email Reminder

[Service]
User=aleph
Type=simple
ExecStart=/home/aleph/.config/remmail.sh
