# remind/.config/systemd/user/r2e.timer
[Unit]
Description=Timer for RSS

[Timer]
RandomizedDelaySec=5
OnBootSec=300s
OnCalendar=daily
Persistent=True
[Install]
WantedBy=timers.target
