# remind/.config/systemd/user/r2e.service
[Unit]
Description=rss2email service

[Service]
User=aleph
Type=simple
ExecStart=/usr/bin/r2e run
